# conf/systemd.service
[Unit]
Description=LobeHub v2 Full Stack (ParadeDB + RustFS + Redis + SearXNG + LobeHub)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=15
TimeoutStartSec=300
TimeoutStopSec=60

WorkingDirectory=__INSTALL_DIR__

ExecStart=/usr/bin/docker-compose \
    --project-name lobehub___APP__ \
    --file __INSTALL_DIR__/docker-compose.yml \
    up \
    --remove-orphans

ExecStop=/usr/bin/docker-compose \
    --project-name lobehub___APP__ \
    down \
    --timeout 30

ExecStopPost=-/usr/bin/docker-compose \
    --project-name lobehub___APP__ \
    down \
    --timeout 10

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
